# vendor/github.com/k8stopologyawareschedwg/deployer/pkg/assets/rte/selinuxinstall.service.template
[Unit]
Description=Install RTE SELinux policy
Before=kubelet.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/semodule -i {{.selinuxPolicyDst}}

[Install]
WantedBy=multi-user.target
